# systemd/verify-remote.service
[Unit]


[Service]
ExecStart=/usr/local/sbin/verify-remote_serv

[Install]
WantedBy=default.target
